# airootfs/usr/lib/systemd/system/crust-update.timer
# Crust OS — periodic update check timer
[Unit]
Description=Periodic Crust OS update check
Documentation=man:crust-update-check(1)

[Timer]
OnBootSec=10min
OnUnitActiveSec=6h
RandomizedDelaySec=30min
Persistent=true

[Install]
WantedBy=timers.target
