# airootfs/usr/lib/systemd/system/crust-recovery-boot.service
# Crust OS — Regenerate snapshot boot entries on every boot
[Unit]
Description=Regenerate systemd-boot entries for BTRFS snapshots
Documentation=man:snapper(8)
ConditionPathExists=/.snapshots
After=local-fs.target

[Service]
Type=oneshot
ExecStart=/usr/local/bin/crust-snap-boot
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target
